# webcam.cron
# make a screenshot
0 10/15 * * * ./var/www/camera/make_screenshot.sh

# send to google drive
